# root/etc/services.d/apache/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/httpd \
  -d /etc/httpd \
  -f conf/httpd.conf \
  -DFOREGROUND
